# support/com.distill.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.distill.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nikpatel/.local/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/Users/nikpatel/Documents/GitHub/distill</string>
        <string>python</string>
        <string>-m</string>
        <string>distill</string>
        <string>run</string>
        <string>--dir</string>
        <string>.</string>
        <string>--output</string>
        <string>./insights</string>
        <string>--use-defaults</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nikpatel/Documents/GitHub/distill</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/nikpatel/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/nikpatel/Library/Logs/distill-daily.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nikpatel/Library/Logs/distill-daily.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
